# docker-proj/systemd-example/bot.service
[Unit]
Description=Run Bot Mon-Fri, 04:00

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run --rm bot-image